# systemd.DISABLED/backtest-automation.service
[Unit]
Description=Backtest Automation Service
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/mnt/e/Trading-Development/2.ブレイクアウト手法プロジェクト
ExecStart=/bin/bash -c 'export DBUS_SESSION_BUS_ADDRESS="unix:path/$XDG_RUNTIME_DIR/bus" && /usr/bin/python3 /mnt/e/Trading-Development/2.ブレイクアウト手法プロジェクト/Scripts/backtest_automation.py'
StandardOutput=journal
StandardError=journal
Environment=PYTHONPATH=/mnt/e/Trading-Development/2.ブレイクアウト手法プロジェクト
Environment=XDG_RUNTIME_DIR=/run/user/1000

[Install]
WantedBy=default.target
